# FRR/supervisord.conf
[supervisord]
nodaemon=true

[program:frr]
command=/usr/lib/frr/docker-start

[program:startup]
command=/etc/rc.local
autorestart=false
priority=1
exitcodes=0
